# .local/service/shell/finish
#!/bin/bash
#exec $HOME/.local/share/go/bin/ags quit -i bar
#exec $HOME/.local/share/go/bin/ags quit -i notif
