# .local/service/udiskie/run
#!/bin/bash
exec udiskie --no-automount
