# .local/service/polkit/run
#!/bin/bash
exec /usr/libexec/hyprpolkitagent
