# .local/service/syncthing/run
#!/bin/bash
export STDATADIR=$HOME/.local/state/syncthing
export STCONFDIR=$HOME/.config/syncthing

exec syncthing
