# .local/service/easyeffects/run
#!/bin/bash
#sv -w 5 check pipewire || exit 1
#exec easyeffects --gapplication-service
sv start pipewire
exec easyeffects -w
