# .local/service/jupyter/run
#!/bin/bash

export PWD=$HOME/Documents/jupyter

exec jupyter-lab
