# .local/service/pipewire/run
#!/bin/bash
exec gentoo-pipewire-launcher restart
