# .local/service/solaar/run
#!/bin/bash
exec solaar -w hide --battery-icons=solaar
